# powertop.service
[Unit]
Description=Startup
After=tlp-init.service
 
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/startup.sh
 
[Install]
WantedBy=multi-user.target
